# messageportd.service.in
[Unit]
Description=Messageport Daemon

[Service]
Type=simple
User=app
SmackProcessLabel=User
KillSignal=SIGTERM
Restart=always
ExecStart=@prefix@/bin/messageportd

[Install]
WantedBy=multi-user.target
